# make_matcher_backend/crontab
0 0 * * * /rails/bin/rails match RAILS_ENV=production >> /rails/log/production.log 2>&1
